# source/devel/pkg-config/profile.d/pkg-config-64.new
#!/bin/sh

if [ -z $PKG_CONFIG_PATH ]; then
  PKG_CONFIG_PATH=/usr/local/lib64/pkgconfig:/usr/lib64/pkgconfig
else
  if ! echo "$PKG_CONFIG_PATH" | grep -q usr/lib64/pkgconfig ; then
    PKG_CONFIG_PATH=$PKG_CONFIG_PATH:/usr/lib64/pkgconfig
  fi
fi

export PKG_CONFIG_PATH

