# source/base/runitscheme/etc/rc.d/rc.local
#!/bin/sh
#
# rc.local:  Local lines (start-up).
#


