# source/base/runitscheme/etc/sv/getty-tty1/run
#!/bin/sh
/sbin/agetty 38400 tty1 linux && /sbin/utmpset -w tty1
