# source/base/runitscheme/etc/sv/getty-tty2/run
#!/bin/sh
/sbin/agetty 38400 tty2 linux && /sbin/utmpset -w tty2
