# upgrades/source/xinetd/xinetd.d/exec
# Begin /etc/xinetd.d/exec

service exec
{
   disable        = yes
   socket_type    = stream
   wait           = no
   user           = 0
   server         = /usr/sbin/in.rexecd
}

# End /etc/xinetd.d/exec
