# upgrades/source/xinetd/xinetd.d/ntalk
# Begin /etc/xinetd.d/ntalk

service ntalk
{
   disable        = yes
   socket_type    = dgram
   wait           = yes
   user           = 0
   server         = /usr/sbin/in.ntalkd
}

# End /etc/xinetd.d/ntalk
