# archive/sysvinit/etc/rc.d/rc.local
#! /bin/sh
#
# /etc/rc.d/rc.local
#
# Local lines (invoked by rc.M).
#


