# source/base/runitscheme/etc/sv/getty-tty3/run
#!/bin/sh
/sbin/agetty 38400 tty3 linux && /sbin/utmpset -w tty3
